# server_config/gunicorn_megano.socket

[Unit]
Description=gunicorn_megano socket

[Socket]
ListenStream=/run/gunicorn_megano.sock

[Install]
WantedBy=sockets.target
